# claude-session-archive-skill/scripts/launchd/com.jrjohn.crs-distill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jrjohn.crs-distill</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jrjohn/claude-archive/distill-backfill.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>Nice</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/jrjohn/claude-archive/distill-backfill.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jrjohn/claude-archive/distill-backfill.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
