# claude-session-archive-skill/templates/crs-build.cron
# /etc/cron.d/crs-build — ingest container session jsonl into shared archive_main every 15 min
# Reads ~/.claude/projects/*/*.jsonl → bge-m3 embed → PG archive_main.
# Mirrors Mac launchd / Linux systemd timer cadence so cross-machine vsearch stays current.
#
# CUSTOMIZE before installing:
#   1. Replace USER / PASSWORD / PG_HOST in CRS_PG_URL.
#   2. Replace PROJECT_PREFIX_LIST in the cron command with this host's owned
#      project prefixes (comma-separated). Without it, embed-missing and
#      ocr-missing pick up ALL unembedded rows globally — including rows from
#      other hosts sharing the same PG → wasted work + "JSONL not found"
#      warnings for foreign-host paths. Examples:
#        bluesea daily-ci-agent: --project-prefix=-root,-workspace-arcana-book
#        Mac launchd wrapper:    --project-prefix=-Users-jrjohn
#      See cloud-deployment.md §"Embed/OCR work-stealing".
# Cron does NOT inherit shell environment, so the env vars below are mandatory.

OLLAMA_HOST=http://ollama:11434
CRS_PG_URL=postgres://USER:PASSWORD@PG_HOST:5432/archive_main?sslmode=require
PATH=/usr/local/bin:/usr/bin:/bin
HOME=/root

# Min Hour DoM Mon DoW User Command
# flock -n: skip this run if previous crs build is still going. Without flock,
# a slow build (Ollama unhealthy, embed backlog) lets cron stack — multiple
# crs build processes hammer Ollama in parallel and trigger the Stopping
# cascade. See cloud-deployment.md "Ollama Stopping cascade" gotcha.
# --workers 2: cap parallel embed requests to 2. crs default is 8, which on a
# CPU-only ARM host with bge-m3 overwhelms Ollama and recreates the cascade
# even with flock + KEEP_ALIVE=-1 + 4G memory all in place. Bump back up only
# if Ollama is on GPU and you've measured headroom (Mac Metal handles 8 fine).
# --project-prefix: see CUSTOMIZE #2 in the header.
*/15 * * * * claude-agent /usr/bin/flock -n /tmp/crs-build.lock /usr/local/bin/crs build --workers 2 --project-prefix=PROJECT_PREFIX_LIST >> /root/claude-archive/crs-build.log 2>&1
